# deploy/vekterm.service
[Unit]
Description=vekterm — USB-DVG/vecterm receiver for the PiTrex
Documentation=https://github.com/anarkiwi/vekterm
# Best-effort wait for the USB serial gadget; Restart=always covers the rest if
# the device appears a little later than the service.
After=dev-ttyGS0.device
Wants=dev-ttyGS0.device

[Service]
Type=simple
ExecStart=/usr/local/bin/vekterm --port /dev/ttyGS0 --baud 2000000
# vekterm exits non-zero if the port is not ready yet; just keep trying.
Restart=always
RestartSec=2
# libpitrex drives the Vectrex via the 6522 VIA over GPIO, which needs hardware
# access; on a dedicated PiTrex appliance the simplest path is to run as root.
User=root
# Give the beam-drawing loop priority so frames stay steady.
Nice=-10

[Install]
WantedBy=multi-user.target
